# infra/production/supervisor/fms-queue.conf
[program:fms-queue]
directory=/var/www/fms-app/current/backend
command=/usr/bin/php artisan queue:work redis --sleep=3 --tries=3 --timeout=120 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/fms/queue-worker.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
stopwaitsecs=3600
environment=APP_ENV="production"
